# examples/systemd/innerwarden-agent.service
[Unit]
Description=Inner Warden - Agent (AI-assisted response)
After=network-online.target innerwarden-sensor.service
Wants=network-online.target
Documentation=https://github.com/InnerWarden/innerwarden

[Service]
Type=simple
User=innerwarden
Group=innerwarden
SupplementaryGroups=systemd-journal
# UMask 0007 keeps new files mode 660 and directories 770. Data in
# /var/lib/innerwarden is security-sensitive (incidents, decisions,
# attacker profiles, dashboard cert/key, graph snapshots). World bits
# must be cleared; group members (sensor runs as root:innerwarden)
# keep read+write for the sensor/agent coordination.
UMask=0007
ExecStart=/usr/local/bin/innerwarden-agent --data-dir /var/lib/innerwarden --config /etc/innerwarden/agent.toml --dashboard --dashboard-bind 0.0.0.0:8787
Restart=on-failure
RestartSec=5

# Environment - API keys, dashboard credentials, Telegram tokens
EnvironmentFile=-/etc/innerwarden/agent.env

# Give it time to flush state and audit trail on shutdown
TimeoutStopSec=15
KillSignal=SIGTERM
SendSIGKILL=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=innerwarden-agent

# SIGHUP triggers config reload (future: hot-reload rules/allowlists)
ExecReload=/bin/kill -HUP $MAINPID

# ── Hardening ──
# NoNewPrivileges omitted: block-ip skills call sudo ufw/iptables/nftables.
# Enable only if using CAP_NET_ADMIN directly or no block-ip skills.
PrivateTmp=yes
ProtectSystem=strict
# /run       - UFW needs /run/ufw.lock
# /etc/ufw   - UFW writes updated rules
# /run/redis  - optional Redis socket
ReadWritePaths=/var/lib/innerwarden /run /etc/ufw
ReadOnlyPaths=/etc/innerwarden
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
