# examples/systemd/innerwarden-sensor.service
[Unit]
Description=Inner Warden - Sensor (host observability)
After=network.target syslog.target
Documentation=https://github.com/InnerWarden/innerwarden

[Service]
Type=simple
User=root
# Group=innerwarden so files created by the sensor are group-owned by
# the agent's user, allowing both sensor (root) and agent (innerwarden)
# to write to incidents-*.jsonl without permission conflicts.
Group=innerwarden
SupplementaryGroups=adm systemd-journal
# UMask 0007 makes new files mode 660 and directories 770. Group write
# (combined with Group=innerwarden above) lets the agent append to files
# created by the sensor. Without this, killchain_inline incidents are
# silently lost with "Permission denied" errors. World bits are cleared
# because /var/lib/innerwarden holds security-sensitive state (incidents,
# decisions, attacker profiles, keys) and must not leak to other users.
UMask=0007
ExecStart=/usr/local/bin/innerwarden-sensor --config /etc/innerwarden/config.toml
Restart=on-failure
RestartSec=5

# Give it time to flush state on shutdown
TimeoutStopSec=10
KillSignal=SIGTERM
SendSIGKILL=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=innerwarden-sensor

# ── Hardening ──
# The sensor runs as root because eBPF requires CAP_BPF + CAP_SYS_ADMIN.
# We lock down everything else to minimize pivot risk if compromised.
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_NET_RAW
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/innerwarden
ReadOnlyPaths=/var/log /etc/innerwarden /proc /sys/fs/bpf
ProtectHome=yes
ProtectKernelTunables=no
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
LockPersonality=yes
MemoryDenyWriteExecute=no

[Install]
WantedBy=multi-user.target
