# examples/systemd/innerwarden-watchdog.service.d/jeprof.conf
# jemalloc tuning drop-in for innerwarden-watchdog.service.
#
# Copy to /etc/systemd/system/innerwarden-watchdog.service.d/jeprof.conf,
# then `sudo systemctl daemon-reload && sudo systemctl restart innerwarden-watchdog`.
#
# Two concerns this drop-in addresses:
#
# 1. Heap-profile readiness. `prof:true` compiles the jemalloc profiler
#    machinery into the running agent so an operator investigating a
#    leak / memory regression can dump a profile without redeploying.
#    `prof_active:false` keeps the profiler DORMANT under normal load;
#    no callsite metadata is collected, no automatic dumps fire.
#    Activate when needed by either:
#      - editing this line to `prof_active:true`, then
#        `sudo systemctl restart innerwarden-watchdog`, or
#      - calling `mallctl("prof.active", true, ...)` from an in-agent
#        signal handler (not currently wired).
#
# 2. Tail-aggressive memory return. `dirty_decay_ms` and `muzzy_decay_ms`
#    at 1000 ms tell jemalloc to release freed pages back to the OS one
#    second after they go quiet, instead of jemalloc's default
#    10-second lazy release. Helps the cgroup avoid holding stale dirty
#    pages on bursty workloads (incident floods). `background_thread`
#    splits the purge work off the request path so the foreground
#    tokio runtime never blocks on madvise.
#
# History: a previous revision shipped `prof_active:true` plus
# `lg_prof_interval:25` (auto-dump every 32 MB allocated). On a busy
# host this generated ~5 GB of dump files and ~50 MB of resident
# callsite metadata in the agent's heap. The 2026-05-21 audit
# (`memory.md` follow-up) caught the configuration drift and pinned the
# defaults below.
[Service]
Environment="MALLOC_CONF=prof:true,prof_active:false,prof_prefix:/var/lib/innerwarden/jeprof/jeprof,lg_prof_sample:19,background_thread:true,dirty_decay_ms:1000,muzzy_decay_ms:1000"
